# debian/ka9q-radio-common.postinst
case "$1" in
  configure)
    deb-systemd-helper enable set_lo_multicast.service >/dev/null || true
    deb-systemd-invoke start set_lo_multicast.service >/dev/null || true
    ;;
esac

#DEBHELPER#
