# debian/ka9q-radio-ft/DEBIAN/postinst
case "$1" in
  configure)
    addgroup --system radio || true
    install -d -m 0775 -o radio -g radio /var/lib/ka9q-radio/ft8
    install -d -m 0775 -o radio -g radio /var/lib/ka9q-radio/ft4
    install -d -m 0775 -o radio -g radio /var/lib/ka9q-radio/wspr
    /usr/bin/set_lo_multicast || true
    ;;
esac

# Automatically added by dh_installtmpfiles/13.24.2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -x "$(command -v systemd-tmpfiles)" ]; then
		systemd-tmpfiles ${DPKG_ROOT:+--root="$DPKG_ROOT"} --create ka9q-radio-ft.conf || true
	fi
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.24.2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if deb-systemd-helper debian-installed 'ft4-decode.service'; then
		# The following line should be removed in trixie or trixie+1
		deb-systemd-helper unmask 'ft4-decode.service' >/dev/null || true

		if deb-systemd-helper --quiet was-enabled 'ft4-decode.service'; then
			# Create new symlinks, if any.
			deb-systemd-helper enable 'ft4-decode.service' >/dev/null || true
		fi
	fi

	# Update the statefile to add new symlinks (if any), which need to be cleaned
	# up on purge. Also remove old symlinks.
	deb-systemd-helper update-state 'ft4-decode.service' >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.24.2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if deb-systemd-helper debian-installed 'ft4-record.service'; then
		# The following line should be removed in trixie or trixie+1
		deb-systemd-helper unmask 'ft4-record.service' >/dev/null || true

		if deb-systemd-helper --quiet was-enabled 'ft4-record.service'; then
			# Create new symlinks, if any.
			deb-systemd-helper enable 'ft4-record.service' >/dev/null || true
		fi
	fi

	# Update the statefile to add new symlinks (if any), which need to be cleaned
	# up on purge. Also remove old symlinks.
	deb-systemd-helper update-state 'ft4-record.service' >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.24.2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if deb-systemd-helper debian-installed 'ft8-decode.service'; then
		# The following line should be removed in trixie or trixie+1
		deb-systemd-helper unmask 'ft8-decode.service' >/dev/null || true

		if deb-systemd-helper --quiet was-enabled 'ft8-decode.service'; then
			# Create new symlinks, if any.
			deb-systemd-helper enable 'ft8-decode.service' >/dev/null || true
		fi
	fi

	# Update the statefile to add new symlinks (if any), which need to be cleaned
	# up on purge. Also remove old symlinks.
	deb-systemd-helper update-state 'ft8-decode.service' >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.24.2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if deb-systemd-helper debian-installed 'ft8-record.service'; then
		# The following line should be removed in trixie or trixie+1
		deb-systemd-helper unmask 'ft8-record.service' >/dev/null || true

		if deb-systemd-helper --quiet was-enabled 'ft8-record.service'; then
			# Create new symlinks, if any.
			deb-systemd-helper enable 'ft8-record.service' >/dev/null || true
		fi
	fi

	# Update the statefile to add new symlinks (if any), which need to be cleaned
	# up on purge. Also remove old symlinks.
	deb-systemd-helper update-state 'ft8-record.service' >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.24.2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if deb-systemd-helper debian-installed 'wspr-decoded.service'; then
		# The following line should be removed in trixie or trixie+1
		deb-systemd-helper unmask 'wspr-decoded.service' >/dev/null || true

		if deb-systemd-helper --quiet was-enabled 'wspr-decoded.service'; then
			# Create new symlinks, if any.
			deb-systemd-helper enable 'wspr-decoded.service' >/dev/null || true
		fi
	fi

	# Update the statefile to add new symlinks (if any), which need to be cleaned
	# up on purge. Also remove old symlinks.
	deb-systemd-helper update-state 'wspr-decoded.service' >/dev/null || true
fi
# End automatically added section

